# ops/systemd/fanmind-backup-retention.timer
[Unit]
Description=FanMind scheduled backup retention

[Timer]
OnCalendar=*-*-* 05:00:00
Persistent=true
Unit=fanmind-backup-retention.service

[Install]
WantedBy=timers.target
