# ops/systemd/fanmind-backup-worker.service
[Unit]
Description=FanMind backup worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=/etc/fanmind-backup/worker.env
EnvironmentFile=-/etc/fanmind-backup/release.env
WorkingDirectory=/usr/local/lib/fanmind-ops
ExecStart=/usr/bin/node /usr/local/lib/fanmind-ops/backup-worker.mjs
Restart=on-failure
RestartSec=30
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/backups/fanmind /tmp
ReadOnlyPaths=/var/www/fanmind/.env.production /etc/nginx /etc/systemd/system /etc/fanmind-backup /home/ubuntu/.pm2/dump.pm2
LockPersonality=true
# MemoryDenyWriteExecute is intentionally omitted: Node.js/V8 JIT on Ubuntu 24.04/Node 24 can require writable executable memory during startup.

[Install]
WantedBy=multi-user.target
